# packaging/debian/prerm
#!/bin/bash
set -e
echo "Removing LockChain ZFS..."
systemctl stop lockchain-zfs.service || true
systemctl disable lockchain-zfs.service || true
