# packaging/systemd/lockchain-key-usb.service
[Unit]
Description=LockChain USB Key Watcher
After=local-fs.target

[Service]
Type=simple
User=root
Group=root
Environment=LOCKCHAIN_CONFIG=/etc/lockchain-zfs.toml
ExecStart=/usr/bin/lockchain-key-usb --config ${LOCKCHAIN_CONFIG}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
