# packaging/systemd/lockchain-zfs.service
[Unit]
Description=LockChain ZFS Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lockchain
Group=lockchain
WorkingDirectory=/var/lib/lockchain
RuntimeDirectory=lockchain
RuntimeDirectoryMode=0750
Environment=LOCKCHAIN_CONFIG=/etc/lockchain-zfs.toml
Environment=LOCKCHAIN_HEALTH_ADDR=127.0.0.1:8787
ExecStart=/usr/bin/lockchain-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
